# examples/launchd/com.email-nurse-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Daily backup of autopilot.yaml. Runs at 03:00 and rotates older copies
  beyond the cap defined in scripts/backup-autopilot.sh.

  Replace HOME_DIR with your absolute home directory before installing.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.email-nurse-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>HOME_DIR/.local/share/email-nurse/current/scripts/backup-autopilot.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>HOME_DIR/Library/Logs/email-nurse-backup.log</string>

    <key>StandardErrorPath</key>
    <string>HOME_DIR/Library/Logs/email-nurse-backup.log</string>
</dict>
</plist>
