# examples/launchd/com.email-nurse-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Optional real-time watcher (KeepAlive while Mail.app is running).
  Most users prefer the scheduled com.email-nurse plist instead; this is
  here for cases where you want continuous polling.

  Replace HOME_DIR with your absolute home directory before installing.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.email-nurse-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>HOME_DIR/.local/bin/email-nurse-watcher.sh</string>
    </array>

    <!-- Only run while Mail.app is running -->
    <key>KeepAlive</key>
    <dict>
        <key>OtherJobEnabled</key>
        <dict>
            <key>com.apple.mail</key>
            <true/>
        </dict>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>HOME_DIR/Library/Logs/email-nurse-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>HOME_DIR/Library/Logs/email-nurse-watcher-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
