# prod-configs/launchd/com.bss.email-nurse-mail-restart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bss.email-nurse-mail-restart</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/bss/code/email-nurse/scripts/mail-restart.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/bss/Library/Logs/email-nurse-mail-restart.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bss/Library/Logs/email-nurse-mail-restart.log</string>
</dict>
</plist>
